# config/fleet.socket
[Socket]
ListenStream=/run/fleet.sock

[Install]
WantedBy=sockets.target
